# server_setup/orderpcbs.conf
[program:orderpcb]
command = /var/www/sites/orderpcbs.com/gunicornstart
user = ec2-user
stdout_logfile = /var/www/sites/orderpcbs.com/logs/gunicorn_supervisor.log
redirect_stderr = true
